# supervisor/jarmilka.conf
[program:jarmilka]
command=python jarmilka.py
user=pi
directory=/application
